# services/LinkyRPiListen.service
[Unit]
Description=LinkyRPi Listen (C version)
After=multi-user.target LinkyRPiGUI.service
Requires=LinkyRPiGUI.service

[Service]
Type=simple
ExecStart=/home/linkyrpi/LinkyRPi/LinkyRPiListen
StandardOutput=append:/var/log/LinkyRPi/LinkyRPiListen.log
StandardError=append:/var/log/LinkyRPi/LinkyRPiListen.err
WorkingDirectory=/home/linkyrpi/LinkyRPi
User=pi
Group=pi
Restart=on-failure
RestartSec=5s
SyslogIdentifier=LinkyRPiListen

# Configuration des permissions et de l'environnement
UMask=0002
Environment=LD_LIBRARY_PATH=/usr/local/lib

# Limites de ressources (optionnel, pour éviter les abus)
MemoryLimit=50M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
